# etc/upstart.conf
description "giraffe"
start on (filesystem)
stop on runlevel [016]

respawn
console log
setuid nobody
setgid nogroup
chdir /opt/app

exec /opt/app/app.sh
